# supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:cron]
command=/usr/sbin/cron -f
user=root
autostart=true

[program:nodejs]
command=/usr/local/bin/npm start
redirect_stderr=true
